# docs/ops/com.einstein.autonomous-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd template. Replace __PROJECT_ROOT__ before installing. See docs/ops/README.md for the install walkthrough. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.einstein.autonomous-loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "__PROJECT_ROOT__" &amp;&amp; uv run python scripts/autonomous_loop.py --max-problems 1 --skip-if-recent 60</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/../mb/autonomous-loop.out.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/../mb/autonomous-loop.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>60</integer>
</dict>
</plist>
