# docs/ops/com.einstein.meta-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd template for the meta-loop (Goal 6 of js/feat/meta-loop).

Fires `scripts/meta_loop.py propose` once per day. Distinct from the inner
autonomous-loop plist — runs much less often, has its own token budget
(mb/logs/meta-loop-budget.md) and audit log (mb/logs/meta-proposals.md).

Replace __PROJECT_ROOT__ with your local checkout's absolute path before
installing. See docs/ops/README.md for the install walkthrough. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.einstein.meta-loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "__PROJECT_ROOT__" &amp;&amp; uv run python scripts/meta_loop.py propose</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <!-- Once per day at 04:00 local. The inner loop runs every 30 min; the
       meta-loop has nothing to do that often — cycle artifacts accumulate
       across the day and one daily proposer pass is plenty. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/../mb/meta-loop.out.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/../mb/meta-loop.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <!-- Set EINSTEIN_META_LOOP=0 to kill-switch all meta-loop proposals at
         the gate chain. Comment out / remove the env var to enable. -->
  </dict>

  <key>ThrottleInterval</key>
  <integer>300</integer>
</dict>
</plist>
